# conky/conky.service
[Unit]
Description=Conky
[Service]
ExecStart=/usr/bin/conky -d
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target